# provision/development/fpm/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:fpm]
command=/usr/bin/php7.0 -R --nodaemonize

[program:sshd]
command=/usr/sbin/sshd -D
